# trackpadctl.cron
#
# Run the trackpadctl program at bootup.
#
@reboot root /usr/bin/trackpadctl >/dev/null 2>&1
